# lab6-mac/app/app.service.j2
[Unit]
Description=app1 daemon
After=network.target

[Service]
Environment=HOSTNAME={{ ansible_fqdn }}
ExecStart=/usr/bin/python3 /root/app.py

[Install]
WantedBy=multi-user.target
